# cronie.ini
[program:cronie]
command=/usr/sbin/crond -f -m off
process_name=%(program_name)s
directory=/tmp
startsecs=5
autostart=true
startretries=3
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
