# memcached.ini
[program:memcached]
command=/usr/bin/memcached -p 11211 -u sogo -m 256 -c 1024 -t 4
user=root
process_name=%(program_name)s
directory=/tmp
startsecs=10
autostart=true
startretries=3
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
