# sogod.ini
[program:sogod]
environment=LD_PRELOAD="/usr/lib/libytnef.so"
command=/usr/bin/sogod -WOPidFile /var/run/sogo/sogo.pid -WONoDetach YES -WOPort 0.0.0.0:20000 -WOLogFile -
process_name=%(program_name)s
directory=/tmp
startsecs=30
autostart=true
startretries=3
autorestart=true
user=sogo
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
